# ansible02/roles/install_hotels_service/files/envoy-api.service
# /etc/systemd/system/envoy-api.service

[Unit]
Description=Envoy Proxy for API Service
After=network.target consul.service
Requires=consul.service

[Service]
ExecStart=/usr/local/bin/envoy -c /etc/consul.d/envoy-api.yaml
Restart=always
RestartSec=5
User=envoy
Group=envoy

[Install]
WantedBy=multi-user.target